# web/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
user=hbnet_web
command=/home/hbnet_web/.local/bin/uwsgi --ini /opt/hbnet_web/wsgi.ini --die-on-term
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
#user=www-data
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
